# code/112.12.125.194(es)/conf/cron/check_pre_ed.cron
*/5 * * * * root /opt/bin/Check_Premium_Edition.py >/dev/null 2>&1
*/5 * * * * root /opt/bin/Check_bjpos.com.py >/dev/null 2>&1
