# code/112.12.125.194(es)/conf/cron/check_junwang_404status.cron
*/5 * * * * root /opt/bin/Check_junwang_404status.py >/dev/null 2>&1
