# code/112.12.125.194(es)/conf/cron/iptables.cron
*/3 * * * * root /usr/bin/iptables_chk.sh > /dev/null 2>&1
