# code/112.12.125.194(es)/conf/cron/salt_agent.cron
*/5 * * * * root /usr/bin/salt.sh >>/var/log/salt.log
